# docker-cleanup.service
[Unit]
Description=Docker Cleanup Service
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/docker-cleanup.sh
User=root
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target 